# contrib/sessions/systemd/user/wayland-wm@margo-session.service.d/10-session-lifecycle.conf
# UWSM's `wayland-wm@.service` template is parameterised by the session
# name (here: `margo-session`, the desktop ID minus `.desktop`). This
# drop-in adds margo-specific tuning without rewriting the template.
#
# Install at:
#   ~/.config/systemd/user/wayland-wm@margo\x2dsession.service.d/10-session-lifecycle.conf
# (the `\x2d` escape is systemd's encoding of `-` inside template
# instance names; keep it literal in the directory name).
#
# After dropping the file:
#   systemctl --user daemon-reload

[Service]
# Compositor log filter (tracing). Picked up by margo and by start-margo
# (start-margo prefers START_MARGO_LOG but falls back to MARGO_LOG when
# the start-specific variable is unset).
Environment=MARGO_LOG=info

# Once the compositor service is up, fan out to a session-wide target
# that other margo-aware user services attach themselves to via
# `WantedBy=margo-session.target` / `PartOf=margo-session.target`.
# Drop the line if you're not running a target-based session model.
ExecStartPost=/usr/bin/systemctl --user start --no-block margo-session.target

StandardOutput=journal
StandardError=journal

# Schedule margo with a small priority bump so the render loop stays
# responsive even when a heavy background user-service (file indexing,
# IDE language server, …) is competing. -15 is the typical "interactive
# compositor" nice value; CPUWeight 500 is double the default 100.
Nice=-15
CPUWeight=500
